# crontab
0 0 * * * echo $(pacman -Qqen) > '/home/sofia/dotfiles/pacman-packages'
0 0 * * * echo $(pacman -Qqem) > '/home/sofia/dotfiles/aur-packages'
